# host/openpd-daemon.service
[Unit]
Description=Data acquisition daemon for OpenPD photodiode amplifier

[Service]
ExecStart=/usr/local/bin/openpd-daemon